# elk_docker_files/containers/supervisord.conf
[supervisord]
nodaemon=true

[program:starter]
command=/start.sh
redirect_stderr=true

[program:kibana]
command=/opt/kibana/bin/kibana
redirect_stderr=true

