# elk_docker_files/logstash_config/supervisor_logstash.conf
[program:logstash]
command=/opt/logstash/bin/logstash -f /config/logstash.conf
stdout_logfile=/var/log/logstash/stdout.log
redirect_stderr=true
