# elk_docker_files/supervisord.conf
[supervisord]
nodaemon=true

[program:starter]
command=/start.sh

[program:kibana]
command=/opt/kibana/bin/kibana
